# root/etc/s6-overlay/s6-rc.d/init-base/up
# This file doesn't do anything, it's just the end of the base image init process
