# root/etc/s6-overlay/s6-rc.d/init-services/up
# This file doesn't do anything, it's just the start of the service init process
