# root/etc/s6-overlay/s6-rc.d/00-legacy/up
ln -s /run/s6/legacy-services /run/s6/services
